# honcho/data/export/upstart/process.conf
start on starting {{app}}-{{name}}
stop on stopping {{app}}-{{name}}
respawn

exec su - {{user}} -s {{shell}} -c 'cd {{app_root}}; export PORT={{port}};{% for var, val in environment.iteritems() %} export {{var|upper}}="{{val}}";{% endfor %} {{ command }} >> {{ log }}/{{ name }}-{{ num }}.log 2>&1'
