# honcho/data/export/supervisord/supervisord.conf
{% for name,command,full_name,num,env in commands %}
[program:{{ full_name }}]
command={{ command }}
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile={{ log }}/{{ name }}-{{ num }}.log
stderr_logfile={{ log }}/{{ name }}-{{ num }}.error.log
user={{ user }}
directory={{ app_root }}
environment={% for key,value in env %}{{ key }}={{ value }}{% if not loop.last %},{% endif %}{% endfor %}{% endfor %}

[group:{{ app }}]
programs={% for name,command,full_name,num,env in commands %}{{ full_name }}{% if not loop.last %},{% endif %}{% endfor %}
