# honcho/export/templates/systemd/master.target
[Unit]
Wants={{ master_wants }}

[Install]
WantedBy=multi-user.target
