# honcho/export/templates/systemd/process_master.target
[Unit]
PartOf={{ app }}.target
Wants={{ process_master_wants }}
